# wallpaper_updater.service
[Unit]
Description=Wallpaper Updater

[Service]
Type=simple
ExecStart=/usr/bin/wallpaper_updater.py

[Install]
WantedBy=multi-user.target
